# automation/com.eric.dailybackup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.osxbackup.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/eric/Dropbox/Development Projects/osx_backup_restore/automation/daily_backup.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <!-- Trigger daily at 15:15 (3:15 PM); the script skips if a backup already completed today to avoid duplicates. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>15</integer>
    <key>Minute</key><integer>15</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/daily_backup.out</string>
  <key>StandardErrorPath</key>
  <string>/tmp/daily_backup.err</string>

</dict>
</plist>
